# systemd/git-deploy-watcher.service
[Unit]
Description=Multi-repo git deploy watcher (SSH + start.sh)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=git-deploy
Group=git-deploy
WorkingDirectory=/var/lib/git-deploy-watcher
EnvironmentFile=-/etc/git-deployer/secrets.env
# Install tree must contain run_watcher.py and git_deploy_watcher/ (see README).
ExecStart=/usr/bin/python3 /opt/git-deploy-watcher/run_watcher.py --config /etc/git-deployer/config.json
Restart=always
RestartSec=10

# Hardening (optional; adjust if start.sh needs more)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
